# roles/infra/systemd_user_timer/templates/ansible-tools-update.timer.j2
[Unit]
Description=Weekly Ansible user-tools update

[Timer]
OnCalendar={{ update_oncalendar }}
Persistent=true

[Install]
WantedBy=default.target
